# edskmgr-support/88-external-disk.rules
#
# /lib/udev/rules.d/88-external-disk.rules
#
# this file is part of the edskmgr package.
# do not edit this file, it will be overwritten on update
ACTION=="add", ENV{DEVTYPE}=="partition", ENV{ID_FS_LABEL}=="?*", ENV{ID_FS_UUID}=="?*", RUN+="/lib/udev/external-disk"
ACTION=="remove", ENV{DEVTYPE}=="partition", ENV{ID_FS_LABEL}=="?*", ENV{ID_FS_UUID}=="?*", RUN+="/lib/udev/external-disk"
ACTION=="add", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", ENV{ID_FS_LABEL}=="?*", ENV{ID_FS_UUID}=="?*", RUN+="/lib/udev/external-disk"
ACTION=="remove", SUBSYSTEM=="block", ENV{DEVTYPE}=="disk", ENV{ID_FS_LABEL}=="?*", ENV{ID_FS_UUID}=="?*", RUN+="/lib/udev/external-disk"
#
# EOF
#
